# ansible/templates/postgresql/detect-rouge-primary.timer.j2
[Unit]
Description=PostgreSQL Split-Brain Detection Timer
Documentation=man:systemd.timer(5)
Requires=detect-rouge-primary.service

[Timer]
# Run every 30 seconds
OnCalendar=*:*:0/30
AccuracySec=1s

# Run immediately after PostgreSQL starts
OnBootSec=60s

# Persistent timer (run missed executions if system was down)
Persistent=false

[Install]
WantedBy=timers.target
