# ansible/files/serve-assets.service
[Service]
ExecStart=/usr/bin/python3 -m http.server 8080
WorkingDirectory=/opt/assets/
